# rtc_egress.service
[Unit]
Description=Agora RTC Egress Server
After=network.target

[Service]
Type=simple
User=egress
Group=egress
WorkingDirectory=/opt/rtc_egress

# Environment variables
Environment=GIN_MODE=release

# Main command to run
ExecStart=/opt/rtc_egress/bin/egress --config /etc/egress/config/egress_config.yaml

# Make sure log directory exists and is writable
ExecStartPre=/bin/mkdir -p /var/log/rtc_egress
ExecStartPre=/bin/chown -R egress:egress /var/log/rtc_egress
ExecStartPre=/bin/chmod 755 /var/log/rtc_egress

# Ensure recordings directory exists and is writable
ExecStartPre=/bin/mkdir -p /recordings
ExecStartPre=/bin/chown -R egress:egress /recordings
ExecStartPre=/bin/chmod 755 /recordings

# Security options
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/recordings /var/log/rtc_egress

# Restart settings
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rtc_egress

[Install]
WantedBy=multi-user.target
